# linux-agent/truenas/menubar-harddrives.timer
[Unit]
Description=Run menubar-harddrives.service periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
